# dist/linux/shub.service
[Unit]
Description=SGX Attestation Hub

[Service]
Type=simple
User=shub
Group=shub
ExecStart=/usr/bin/shub run
ExecReload=/bin/kill -s HUP $MAINPID
TimeoutStartSec=0
Restart=on-failure
PermissionsStartOnly=true
RuntimeDirectory=shub
RuntimeDirectoryMode=0775

[Install]
WantedBy=multi-user.target
